# crontab
*       *       *       *       *       /sbin/pfctl -t abusers -Te 60
*       *       *       *       *       /sbin/pfctl -t uberabusers -Te 360